# nextcloud/etc/crontabs/nextcloud
*/5 * * * * php -f "$(occ config:system:get server_root)/cron.php"
