# nextcloud/etc/crontabs/nextcloud
*/5 * * * * php82 -f /usr/share/webapps/nextcloud/cron.php --define apc.enable_cli=1
